# ansible/roles/core/templates/profile.d/init-env.sh.j2
#! /usr/bin/env bash

export GHOST_AWS_REGION={{ aws_region }}
export GHOST_BUILD_VERSION={{ build_version }}
export GHOST_COMMIT_HASH={{ commit_hash }}
export GHOST_ENVIRONMENT={{ environment_ }}
